# code/dashboard-carousel/image-uploader/config/upstart/image_uploader.conf
description "start and stop the <Your App> node.js server"
version "1.0"
author "ogolan@atti.com"

start on filesystem and started networking
stop on shutdown
respawn limit 99 5

env NODE_ENV=production
env PATH=/home/you/node/bin
chdir /home/you/projects/your_app
exec su -c 'node server.js' www-data  >> /var/log/yourapp.log 2>&1

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/root"

    exec sudo -u node2dm /usr/local/bin/node /srv/node2dm/node2dm.js >> /var/log/node.log 2>&1
end script


